# agent/wol-agent.service.example
[Unit]
Description=WoL Agent — registra IP e serviços no painel
After=network-online.target
Wants=network-online.target

[Service]
User=andre-reis
WorkingDirectory=/home/andre-reis/Documents/wol-rpi-k8s/agent
ExecStart=/home/andre-reis/Documents/wol-rpi-k8s/agent/venv/bin/python agent.py
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
